# git-auto-pull/systemd/git-auto-pull.timer
[Unit]
Description=Hourly git auto-pull timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
